# .config/systemd/user/backup.service
[Unit]
Description=daily backup

[Service]
Type=simple
ExecStart=/home/ajgrf/bin/backup
